# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:citedby_api]
command=gunicorn --paste production.ini -w 4 -b 0.0.0.0

[program:citedby_thriftserver]
command=citedby_thriftserver --port 11620 --host 0.0.0.0